# sv/dashing/run
#!/bin/sh

#sv start socklog-unix || exit 1

#cd /dashboard
cd /$DASHBOARD

#exec 2>&1
exec dashing start -d
